# etc/uci-defaults/96-wwan-modemmanager
#!/bin/sh
# Switch the `wwan` interface from `proto mbim` to `proto modemmanager` once
# on first boot (or first sysupgrade after the MM swap), and seed the APN /
# auth / pdptype fields ModemManager expects.
#
# Why: the umbim/mbim-watchdog combo we used to ship couldn't reliably catch
# carrier-side silent disconnects on the RM520N (modem only emits
# DEACTIVATING for some teardowns and nothing at all for idle-timer drops),
# so the data plane got wedged whenever the carrier renegotiated. MM has
# proper signal-loss / RAT-change / disconnect handling baked in.
#
# Idempotent: only acts if proto != modemmanager. Locating the modem's sysfs
# path is done at runtime because PCIe MHI doesn't give us a stable UCI key
# for it — we walk /sys/class/wwan to find the wwan0mbim0 control interface
# and resolve the parent device that MM expects to see in the `device`
# field.

. /lib/functions.sh

board=$(board_name)
case "$board" in
    glinet,gl-x3000|glinet,gl-xe3000) ;;
    *) exit 0 ;;
esac

[ -f /etc/config/network ] || exit 0
[ "$(uci -q get network.wwan.proto)" = "modemmanager" ] && exit 0

# Resolve the sysfs path of the MBIM control device. Try MHI's wwan subsystem
# first (where the RM520N lives on the GL-X3000) and fall back to the legacy
# usbmisc cdc-wdm0 path so the same uci-defaults works on USB-only modems.
device=
for candidate in \
    /sys/class/wwan/wwan0mbim0/device \
    /sys/class/usbmisc/cdc-wdm0/device; do
    [ -e "$candidate" ] || continue
    device=$(readlink -f "$candidate")
    break
done

[ -n "$device" ] || {
    logger -t wwan-mm-defaults "no MBIM control device found, leaving network.wwan untouched"
    exit 0
}

uci -q batch <<EOF
set network.wwan.proto='modemmanager'
set network.wwan.device='$device'
set network.wwan.apn='wap.tim.it'
set network.wwan.auth='none'
set network.wwan.iptype='ipv4v6'
delete network.wwan.pdptype
delete network.wwan.allow_partner
delete network.wwan.allow_roaming
EOF
uci commit network

logger -t wwan-mm-defaults "wwan switched to proto=modemmanager device=$device"

exit 0
